# produccion/grua_8/servicios/grua8_64_descarga.service
[Unit]
Description=Descarga imagenes camara 64 grua 8 
After=network.target

[Service]
ExecStart=python3 captura_gstreamer.py config/camara_64.ini 
User=safetymind
WorkingDirectory=/data/pipelines-grua-apt/captura
Restart=always  

[Install]
WantedBy=multi-user.target